# scripts/com.shortsfactory.manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.EN.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shortsfactory.manager</string>
    <key>ProgramArguments</key>
    <array>
        <!-- 絶対パスに書き換えてください -->
        <string>/Users/motista/Desktop/antigravity/modular-open-claw/target/release/shorts-factory</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/motista/Desktop/antigravity/modular-open-claw</string>
    <key>StandardOutPath</key>
    <string>/Users/motista/Desktop/antigravity/modular-open-claw/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/motista/Desktop/antigravity/modular-open-claw/logs/stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>
</dict>
</plist>
