# scripts/com.aiome.watchtower.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aiome.watchtower</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/motista/Desktop/antigravity/modular-open-claw/target/release/watchtower</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/motista/Desktop/antigravity/modular-open-claw</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/aiome.watchtower.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/aiome.watchtower.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
